# launchd/com.oscargiller.resy-sniper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oscargiller.resy-sniper.macmini</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/oscar/Projects/ResSniper/.venv/bin/python</string>
        <string>watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/oscar/Projects/ResSniper</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/oscar/Projects/ResSniper/logs/watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/oscar/Projects/ResSniper/logs/watcher.err</string>
</dict>
</plist>
